# deploy/bookshelf-api.service
[Unit]
Description=Bookshelf API (FastAPI + uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bookshelf-api
EnvironmentFile=/etc/bookshelf.env
ExecStart=/opt/bookshelf-api/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
